# redis/redis.ini
[program:redis]
command=/usr/bin/redis-server /etc/redis.conf
directory=/usr/lib/check_mk_agent/push
stdout_logfile=/var/log/redis.log
autostart=true
autorestart=true
redirect_stderr=true

