# redis/redis_sentinel.ini
[program:redis-sentinel]
command=/usr/bin/redis-sentinel /root/bimax-counters/redis/redis-sentinel.conf
directory=/root/bimax-counters/redis/
stdout_logfile=/var/log/redis_sentinel.log
autostart=true
autorestart=true
redirect_stderr=true
